# assets/debian.postinst
#!/bin/sh

systemctl --global enable swaypanion